# scripts/monitoring/aicon-alerting.service
[Unit]
Description=AICON Social Media Analysis Alerting System
Documentation=https://github.com/yourusername/aicon
After=network.target

[Service]
Type=simple
User=aicon
WorkingDirectory=/opt/aicon
EnvironmentFile=/opt/aicon/.env
ExecStart=/usr/bin/node /opt/aicon/scripts/monitoring/alerting.js
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aicon/logs

# Resource limits
MemoryLimit=512M
CPUQuota=25%

# Logging
StandardOutput=append:/opt/aicon/logs/alerting.log
StandardError=append:/opt/aicon/logs/alerting-error.log

[Install]
WantedBy=multi-user.target